# init/msm@.service
[Unit]
Description=Minecraft Server Manager - %i
After=network.target local-fs.target
Wants=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/msm %i start
ExecStop=/usr/local/bin/msm %i stop
ExecReload=/usr/local/bin/msm %i restart

# Allow time for graceful shutdown (player warning + world save)
# 10s player warning + 30s stop wait + buffer for large worlds
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
